# conf/ai.service
[Unit]
Description=Pandora AI service
After=network.target

[Service]
Type=simple
User=__APP__
Group=__APP__
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=__ENV_FILE__
ExecStart=__INSTALL_DIR__/ai/venv/bin/python ai/app.py
Restart=always
RestartSec=5
StandardOutput=append:__LOG_DIR__/ai.log
StandardError=append:__LOG_DIR__/ai.error.log

[Install]
WantedBy=multi-user.target
